# scripts/utils/docker-runner.conf
start on started tty1
respawn
script
    cd /home/ubuntu
    bash ./loop.sh
end script